# inventario/infrastructure/systemd/lunari-inventory-api.service
[Unit]
Description=LUNARi Inventory API Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/usr/bin/java -jar /home/ec2-user/lunari-inventory-api-0.0.1-SNAPSHOT.jar --spring.profiles.active=prod --DB_HOST=your-db.neon.tech --DB_PORT=5432 --DB_NAME=lunari_inventory --DB_USER=your_user --DB_PASSWORD=your_password
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
